# scripts/einkframe.service
[Unit]
Description=E-Ink Photo Frame
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/myungs/eink_frame
ExecStart=/home/myungs/eink_frame/venv/bin/python /home/myungs/eink_frame/run_production.py
Restart=on-failure
RestartSec=10

# Environment
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
